# gitlab-ci/iopsys-supervisord.conf
[program:ubusd]
command=/bin/bash -c "/usr/sbin/ubusd"

[program:uspd]
command=/bin/bash -c "/usr/bin/valgrind --xml=yes --xml-file=memory-report.xml --leak-check=full --show-reachable=yes --show-leak-kinds=all --errors-for-leak-kinds=all --error-exitcode=1 --track-origins=yes /builds/iopsys/uspd/uspd"

